# docker/supervisord.conf
; ============================================================
; Supervisor programs for the production FrankenPHP container.
; Global supervisord config lives at /etc/supervisor/supervisord.conf
; (shipped by the supervisor package). This file is included via
; [include] files = /etc/supervisor/conf.d/*.conf
; ============================================================

; --------- 1. FrankenPHP (web server) ---------
[program:frankenphp]
command=frankenphp run --config /etc/caddy/Caddyfile --adapter caddyfile
autostart=true
autorestart=true
startretries=5
user=www-data
stopsignal=SIGTERM
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; --------- 2. Laravel Queue Worker ---------
[program:queue-worker]
command=php /app/artisan queue:work --tries=3 --timeout=90 --sleep=3 --max-time=3600
autostart=true
autorestart=true
startretries=5
user=www-data
stopsignal=SIGTERM
stopwaitsecs=100
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

; --------- 3. Laravel Scheduler (long-running) ---------
[program:scheduler]
command=php /app/artisan schedule:work
autostart=true
autorestart=true
startretries=5
user=www-data
stopsignal=SIGTERM
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30

; --------- 4. Inertia SSR Server (Node.js) ---------
; `php artisan inertia:start-ssr` akan menjalankan `node bootstrap/ssr/ssr.js`
; dan listen di port yang dikonfigurasi di `config/inertia.php` (default 13714).
; FrankenPHP akan memanggil endpoint tersebut untuk pre-render setiap halaman.
[program:inertia-ssr]
command=php /app/artisan inertia:start-ssr
autostart=true
autorestart=true
startretries=5
user=www-data
stopsignal=SIGTERM
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=5
